# misc/koto-faucet.service
[Unit]
Description=uWSGI Emperor
After=syslog.target

[Service]
# edit all 'username' to your user
ExecStart=/home/username/koto-faucet/start_uwsgi.sh
WorkingDirectory=/home/username/koto-faucet
Restart=always
KillSignal=SIGQUIT
Type=notify
StandardError=syslog
NotifyAccess=all
User=username
Group=username
PIDFile=/home/username/koto-faucet/koto-faucet.pid

[Install]
WantedBy=multi-user.target
